# crontab
*/10 7-9,19-21 * * * export PATH=/usr/local/bin:$PATH; cd /app; npm start >> /var/log/cron.log 2>&1
